# docker/cronjob/crontab
*/5 * * * * php /var/www/symf_project/bin/console newsparser
